# src/temperature_logger.service
[Unit]
Description=Temperature Logger

[Service]
ExecStart=/usr/bin/temperature_logger

[Install]
WantedBy=default.target
